# config/audia.service
[Unit]
Description=Audia - Kiosk server de cribado fonológico
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/audia-server
Environment="FLASK_ENV=production"
Environment="PORT=5000"
ExecStart=/home/pi/audia-server/venv/bin/python app.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
